# nix/packages/swhkd/swhkd.service
[Unit]
Description=swhkd hotkey daemon
BindsTo=default.target

[Service]
Type=simple
Restart=always
ExecStart=@out@

[Install]
WantedBy=default.target
